# crontab
* * * * * '.heroku/php/bin/php' 'artisan' inspire

# An empty line is required at the end of this file for a valid cron file.
